# roles/nas/files/backup-script.service
[Unit]
Description=General Backup Task

[Service]
ExecStart=/media/data/scripts/backup_general.sh