# roles/nas/files/media-backup.automount
[Unit]
Description=Automount Backup

[Automount]
Where=/media/backup

[Install]
WantedBy=multi-user.target
